# conf/supervisord.conf
[program:sender]
command = /home/jumpserver/python2/bin/gunicorn -c ../../conf/gun.py server:app
directory = /home/jumpserver/sender/bin
autostart = true
autorestart = true
